# config/services.d/proxy-hooker/run
#!/bin/execlineb -P
/usr/bin/proxy-hooker